# software/distro/build/base-os/networking/etc/systemd/system/autohotspot.service
[Unit]
Description=Automatic generation of a wireless AP (wifi hotspot) when no known wifi networks are found
After=multi-user.target
[Service]
Type=oneshot
ExecStart=/home/pi/.local/bin/autohotspot.sh
[Install]
WantedBy=multi-user.target
