# software/distro/build/base-os/networking/etc/systemd/system/autohotspot.service
[Unit]
Description=Automatically generates a Hotspot when a valid SSID is not in range
After=multi-user.target
[Service]
Type=oneshot
ExecStart=/home/pi/.local/bin/autohotspot.sh
[Install]
WantedBy=multi-user.target
