# software/distro/build/base-os/networking/etc/systemd/system/autohotspot.timer
[Unit]
Description=Run the autohotspot periodically in case a recognized wifi network comes in range

[Timer]
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
