# software/distro/build/base-os/networking/etc/systemd/system/autohotspot.timer
[Unit]
Description=Run the autohotspot service periodically in case a known wifi network becomes available

[Timer]
OnUnitActiveSec=2min
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target
