# software/distro/build/base-os/networking/etc/systemd/system/first-boot-ssh-keys.service
[Unit]
Description=Regenerates OpenSSH host keys
Before=sshd.service
[Service]
Type=oneshot
ExecStartPre=systemctl disable first-boot-ssh-keys.service
ExecStartPre=rm -f /etc/ssh/ssh_host_*
ExecStart=dpkg-reconfigure --frontend=noninteractive openssh-server
[Install]
WantedBy=sshd.service
