# software/distro/build/planktoscope-app-env/forklift/etc/systemd/system/first-boot-forklift-apply.service
[Unit]
Description=Deploys Forklift pallet
Wants=docker.service
After=docker.service
[Service]
Type=oneshot
ExecStartPre=systemctl disable first-boot-forklift-deploy.service
ExecStart=/home/pi/.local/bin/forklift --workspace /home/pi/.forklift plt apply
[Install]
WantedBy=multi-user.target
