# software/distro/build/planktoscope-app-env/forklift/etc/systemd/system/first-boot-forklift-deploy.service
[Unit]
Description=Deploys forklift environment
Wants=first-boot-docker-swarm.service
After=first-boot-docker-swarm.service
[Service]
Type=oneshot
ExecStartPre=systemctl disable first-boot-forklift-deploy.service
ExecStart=/home/pi/.local/bin/forklift --workspace /home/pi/.forklift env apply
[Install]
WantedBy=multi-user.target
