# software/distro/setup/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.update-cockpit-origins-machine-name.service
[Unit]
Description=Update Cockpit origins based on machine name generated from the Raspberry Pi's serial number
Before=network-pre.service
Wants=network-pre.target

[Service]
Type=oneshot
# Note: `exec` is needed for `$(/usr/bin/forklift ...)` to be evaluated
ExecStart=/usr/bin/bash -c "exec $(/usr/bin/forklift --workspace /home/pi plt locate-depl-pkg host/planktoscope/machine-name)/cockpit-origins/update.sh"
Restart=no

[Install]
WantedBy=multi-user.target
