# software/distro/setup/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.update-ssid-machine-name.service
[Unit]
Description=Update hostapd SSID based on machine name generated from the Raspberry Pi's serial number
Before=hostapd.service

[Service]
Type=oneshot
ExecStart=/home/pi/.local/bin/update-ssid-machine-name.sh
Restart=no

[Install]
WantedBy=hostapd.service
