# software/distro/setup/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.update-cockpit-origins-machine-name.service
[Unit]
Description=Update Cockpit origins based on machine name generated from the Raspberry Pi's serial number
Before=network-pre.service
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart="$(/home/pi/.local/bin/forklift --workspace /home/pi plt locate-depl-pkg host/planktoscope/machine-name)/update-cockpit-origins.sh"
Restart=no

[Install]
WantedBy=multi-user.target
