# software/distro/setup/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.update-machine-name.service
[Unit]
Description=Update a file with the machine name generated from the Raspberry Pi's serial number

[Service]
Type=oneshot
ExecStart="$(/usr/bin/forklift --workspace /home/pi plt locate-depl-pkg host/planktoscope/machine-name)/update-machine-name.sh"
Restart=no

[Install]
WantedBy=multi-user.target
