# controller/light/planktoscope-org.controller.light.service
[Unit]
Wants=mosquitto.service
After=mosquitto.service

[Service]
Type=simple
Environment=HOME=/home/pi
WorkingDirectory=/home/pi/PlanktoScope/controller
ExecStart=/usr/local/bin/uv run python -m light.main
User=pi
Group=pi
Restart=on-failure

[Install]
WantedBy=multi-user.target
