# software/distro/setup/base-os/tools/etc/systemd/system/docker.service.d/override.conf
[Unit]
# Don't wait for network-online.target, since that takes a long time
After=docker.socket firewalld.service containerd.service time-set.target
