# software/distro/build/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.avahi-alias-ether-planktoscope.local.service
[Unit]
Description=Publish ether-planktoscope.local as alias for 192.168.5.1 via mDNS

[Service]
Type=simple
ExecStart=/bin/bash -c "/usr/bin/avahi-publish -a -R ether-planktoscope.local 192.168.5.1"

[Install]
WantedBy=multi-user.target
