# software/distro/build/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.update-hosts-machine-name.service
[Unit]
Description=Update dnsmasq hosts based on machine name generated from the Raspberry Pi's serial number
Before=dnsmasq.service

[Service]
Type=oneshot
ExecStart=/home/pi/.local/bin/update-hosts-machine-name.sh
Restart=no

[Install]
WantedBy=dnsmasq.service
