# software/distro/build/planktoscope-app-env/networking/etc/systemd/system/planktoscope-org.avahi-alias-planktoscope.local.service
[Unit]
Description=Publish planktoscope.local as alias for 192.168.4.1 via mDNS

[Service]
Type=simple
ExecStart=/bin/bash -c "/usr/bin/avahi-publish -a -R planktoscope.local 192.168.4.1"

[Install]
WantedBy=multi-user.target
