# software/distro/setup/base-os/forklift/usr/lib/systemd/system/overlay-fs.target
[Unit]
Description=Local File Systems with Overlays
DefaultDependencies=no
Conflicts=shutdown.target
After=local-fs.target
OnFailure=emergency.target
OnFailureMode=replace-irreversibly
