# os/raspberry/planktoscope-org.firstboot.service
[Unit]
Description=PlanktoScope first boot service
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=/home/pi/PlanktoScope/os/bootloader/firstboot.sh

[Install]
WantedBy=sysinit.target
