# software/distro/setup/planktoscope-app-env/gpio-stepper-init/etc/systemd/system/planktoscope-org.init-gpio-steppers.service
[Unit]
Description=Release GPIO stepper motors on startup
DefaultDependencies=false

[Service]
Type=oneshot
ExecStart=/home/pi/.local/bin/release-gpio-steppers.sh
Restart=no

[Install]
WantedBy=sysinit.target
